# config/deploy/services/puma_IRM_env.service.erb
[Unit]
Description=Puma HTTP Server for IRM (<%= fetch(:stage) %>)
After=network.target


[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/IRM/current
Environment="RAILS-ENV=<%= fetch(:stage) %>"

# Support older bundler versions where file descriptors weren't kept
# See https://github.com/rubygems/rubygems/issues/3254
ExecStart=/home/ubuntu/.rvm/bin/rvm default do bundle exec --keep-file-descriptors puma -C /home/ubuntu/IRM/shared/puma.rb
ExecReload=/bin/kill -USR1 $MAINPID
StandardOutput=append:/home/ubuntu/IRM/current/log/puma.error.log
StandardError=append:/home/ubuntu/IRM/current/log/puma.access.log




Restart=always
RestartSec=1

SyslogIdentifier=puma

[Install]
WantedBy=multi-user.target